# systemd/wireguard-log-connections.service
[Unit]
Description=WireGuard connection audit logger
Documentation=file:///etc/wireguard/scripts/log-connections.sh
After=wg-quick.target

[Service]
Type=oneshot
ExecStart=/etc/wireguard/scripts/log-connections.sh

# Hardening
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/wireguard-connections
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
