# systemd/wireguard-healthcheck.service
[Unit]
Description=WireGuard healthcheck (service + interface + addresses + firewall)
Documentation=file:///home/wireguard-scripts/healthcheck.sh
After=wg-quick.target

[Service]
Type=oneshot
# Optional: set WG_PING_TARGET=<server-tunnel-IP> in this file to enable the
# upstream reachability check (restart the tunnel when the server stops
# answering). The leading '-' makes the file optional — absent = structural
# checks only, the original behavior.
EnvironmentFile=-/etc/wireguard/healthcheck.env
ExecStart=/home/wireguard-scripts/healthcheck.sh --restart

# Minimal hardening. Less restrictive than the connection-log unit because
# this script may need to `systemctl restart wg-quick@<iface>` and inspect
# firewall state via nft/iptables/firewall-cmd — locking it down breaks
# auto-recovery.
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
RestrictRealtime=yes
