# systemd/wireguard-healthcheck.timer
[Unit]
Description=Run WireGuard healthcheck every minute (restart ~4 min into an outage)

[Timer]
OnBootSec=1min
OnUnitActiveSec=60s
RandomizedDelaySec=10
Persistent=true

[Install]
WantedBy=timers.target
